# dotfiles/systemd/.config/systemd/user/checkup.service
[Unit]
Description=Check available updates

[Service]
Type=oneshot
ExecStart=/home/siavash/.local/bin/services/checkup

[Install]
WantedBy=graphical.target
