# home/dotfiles/systemd/user/newsup.timer
[Unit]
Description=Check available news every 1 hour

[Timer]
OnActiveSec=10m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
